# deb/install/postinst
#!/bin/bash

case "${1}" in
  configure)
    # enable the service to run on system shutdown
    systemctl enable waggle-node-hostname.service
  ;;
esac
