# deb/install/prerm
#!/bin/bash

systemctl disable waggle-node-hostname.service
