# ROOTFS/etc/systemd/system/waggle-node-hostname.service
[Unit]
Description=Waggle Hostname Service
Before=network-pre.target
Wants=network-pre.target
ConditionPathExists=/etc/waggle/node-id

[Service]
ExecStart=/usr/bin/waggle_node_hostname.py -n /etc/waggle/node-id
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=default.target
